# conf/systemd.service
[Unit]
Description=Movim daemon
After=nginx.service network.target local-fs.target postgresql.service

[Service]
Type=simple
User=www-data
Environment=PUBLIC_URL=https://__DOMAIN__/__PATH__
Environment=WS_PORT=__PORT__
ExecStart=/usr/bin/php daemon.php start --url=${PUBLIC_URL} --port=${WS_PORT}
WorkingDirectory=__FINALPATH__
StandardOutput=syslog
SyslogIdentifier=__APP__
PIDFile=/run/movim.pid
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
