# conf/systemd.service
[Unit]
Description=Movim daemon
After=nginx.service network.target local-fs.target postgresql.service

[Service]
Type=simple
User=www-data
ExecStart=/usr/bin/php daemon.php start --interface=127.0.0.1 --url=https://__DOMAIN____PATH__ --port=__PORT__
WorkingDirectory=__FINALPATH__
StandardOutput=syslog
SyslogIdentifier=__APP__
PIDFile=/run/movim.pid
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
